# control-plane/akigura.service
[Unit]
Description=AkiGura SaaS Control Plane
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/AkiGura/control-plane
ExecStart=/home/exedev/AkiGura/control-plane/akigura-srv -listen :8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
